# deploy_tools/gunicorn-upstart.template.conf
cription "Gunicorn server for mixees"

start on net-device-up
stop on shutdown

respawn

setuid ubuntu
chdir /home/ubuntu/sites/mixees/source

exec ../virtualenv/bin/gunicorn mixees.wsgi:application
